# launchd-plists/com.agentcomms.poller.builder.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.agentcomms.poller.builder</string>
  <key>ProgramArguments</key>
  <array>
    <string>${HOME}/.openclaw/workspace-builder/cron/agent-poller.sh</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>AGENT_WORKSPACE</key>
    <string>${HOME}/.openclaw/workspace-builder</string>
    <key>AGENT_LOG_DIR</key>
    <string>${HOME}/.openclaw/workspace-builder/logs</string>
    <key>AGENT_IDLE_TIMEOUT</key>
    <string>300</string>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin:${HOME}/.npm-global/bin:${HOME}/Library/pnpm</string>
    <key>HOME</key>
    <string>${HOME}</string>
  </dict>
  <key>StartInterval</key>
  <integer>30</integer>
  <key>StandardOutPath</key>
  <string>${HOME}/.openclaw/workspace-builder/logs/poller.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>${HOME}/.openclaw/workspace-builder/logs/poller.stderr.log</string>
  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>
